# Tproxy/etc/init.d/iptables
#!/bin/sh /etc/rc.common
# OpenWrt init script for TPROXY management (iptables version)
# Place this file at /etc/init.d/tproxy and make it executable

START=99
STOP=15
USE_PROCD=1 # 启用 procd 管理

RESERVED_IPS="
10.0.0.0/8
17.0.0.0/8
100.64.0.0/10
127.0.0.0/8
169.254.0.0/16
172.16.0.0/12
192.0.0.0/24
192.168.0.0/16
224.0.0.0/4
240.0.0.0/4
255.255.255.255/32
115.25.0.0/18
115.25.212.0/22
123.124.147.0/24
202.204.48.0/20
211.144.7.0/24
222.28.32.0/19
222.28.64.0/20
222.199.192.0/19
"

start_service() {
    echo "Starting TPROXY rules (iptables)..."

    # 检查规则是否已存在
    if iptables -t mangle -L PREROUTING -n | grep -q 'TPROXY_RULE'; then
        echo "TPROXY rules are already running."
        return 0
    fi

    # 创建 ipset
    ipset create RESERVED_IP hash:net || ipset flush RESERVED_IP
    for ip in $RESERVED_IPS; do
        ipset add RESERVED_IP "$ip"
    done

    # 创建 TPROXY 链
    iptables -t mangle -N TPROXY_RULE

    # 规则: 跳过保留地址
    iptables -t mangle -A TPROXY_RULE -m set --match-set RESERVED_IP dst -j RETURN

    # 规则: TPROXY 转发
    iptables -t mangle -A TPROXY_RULE -p tcp -j TPROXY --on-port 12345 --on-ip 127.0.0.1 --tproxy-mark 1
    iptables -t mangle -A TPROXY_RULE -p udp -j TPROXY --on-port 12345 --on-ip 127.0.0.1 --tproxy-mark 1

    # 从 PREROUTING 跳转到 TPROXY 链
    iptables -t mangle -A PREROUTING -j TPROXY_RULE

    ulimit -n 1000000  # 增加文件句柄限制

    # 配置 ip rule 和路由表
    ip rule add fwmark 1 table 100
    ip route add local 0.0.0.0/0 dev lo table 100
}

stop_service() {
    echo "Stopping TPROXY rules (iptables)..."

    ip rule del fwmark 1 table 100 || true
    ip route del local 0.0.0.0/0 dev lo table 100 || true

    # 清理 iptables 规则
    iptables -t mangle -F TPROXY_RULE || true
    iptables -t mangle -D PREROUTING -j TPROXY_RULE || true
    iptables -t mangle -X TPROXY_RULE || true

    # 销毁 ipset
    ipset destroy RESERVED_IP || true
}

restart_service() {
    stop_service
    sleep 1
    start_service
}